# deploy/bot-fractal.service
[Unit]
Description=Bot Fractal Mind (agente_mente.py)
After=network-online.target
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/fractal-mind
EnvironmentFile=/home/ubuntu/fractal-mind/.env
ExecStart=/home/ubuntu/bin/micromamba run -n bot-env python /home/ubuntu/fractal-mind/scripts/agente_mente.py
Restart=always
RestartSec=10
StartLimitInterval=0
StartLimitBurst=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot-fractal

[Install]
WantedBy=multi-user.target
