# deploy/bot-air.service
[Unit]
Description=Bot AIR - Videos (agente_air.py)
After=network-online.target
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/fractal-mind
EnvironmentFile=/home/ubuntu/fractal-mind/.env
ExecStart=/home/ubuntu/bin/micromamba run -n bot-env python /home/ubuntu/fractal-mind/proyectos/air-bot/agente_air.py
Restart=always
RestartSec=10
StartLimitInterval=0
StartLimitBurst=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot-air

[Install]
WantedBy=multi-user.target
